# infra/systemd/yolo-vlm-backend.service
[Unit]
Description=YOLO-LLM VLM Inference Backend (llama.cpp)
After=network.target

[Service]
Type=simple
User=ray
WorkingDirectory=/home/ray/developer/ai/YOLO_LLM
ExecStart=/home/ray/developer/ai/YOLO_LLM/scripts/run.sh vlm-backend

Restart=on-failure
RestartSec=10
WatchdogSec=120

[Install]
WantedBy=multi-user.target
