# infra/systemd/yolo-perception.service
[Unit]
Description=YOLO-LLM Perception (Fast Path)
After=network.target mediamtx.service
Requires=mediamtx.service

[Service]
Type=simple
User=ray
WorkingDirectory=/home/ray/developer/ai/YOLO_LLM
ExecStart=/home/ray/developer/ai/YOLO_LLM/scripts/run.sh perception

Restart=always
RestartSec=5
StartLimitIntervalSec=0
WatchdogSec=30

[Install]
WantedBy=multi-user.target
